# crontab
# keep pppoe connecting
* * * * * /bin/sh /etc/cron.d/pppoe-keep_watchdog.sh
